# 4-server_unit/resource/ap/etc/cron.d/log-delete-cron
30 21 * * * root find -L /fj/dc-core/log -name '*.log*' -mtime +31 -exec /bin/rm -f {} \; > /dev/null 2>&1
30 21 * * * root find -L /fj/dc-engine/log -name '*.log*' -mtime +31 -exec /bin/rm -f {} \; > /dev/null 2>&1
